# mydir/files/etc/uci-defaults/99-disable-privoxy
#!/bin/sh
# 禁用 Privoxy 开机自启
/etc/init.d/privoxy disable
exit 0
