# mydir/files/etc/uci-defaults/99-disable-unbound
#!/bin/sh
# 禁用 unbound 开机启动
/etc/init.d/unbound disable
exit 0
